# profile.d/00_basic_variables
# -*- mode: sh; -*-
# vim: set filetype=sh :

# sets general variables if they don't get set elsewhere
: ${USER:=$(whoami)}
: ${HOST:=$(hostname -f)}
: ${HOME:=$(finger $USER | awk '/^Directory/ {print $2 }')}
